# network-config-applier.service
[Unit]
Description=Network Configuration Applier
Documentation=man:systemd.service(5)
After=network.target NetworkManager.service
Wants=network-online.target
ConditionPathExists=/opt/gateway/pending_network_config

[Service]
Type=oneshot
ExecStart=/opt/gateway/network_config_applier.sh
RemainAfterExit=no
User=root
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log /opt/gateway /etc/network /etc/systemd

[Install]
WantedBy=multi-user.target